# git-fetch-daily.cron
#!/bin/sh

export PATH=/srv/git/tools:${PATH}

git-fetch-daily /srv/git 2>&1 | tee /srv/git/git-fetch-daily.log
