# git-fetch-hourly.cron
#!/bin/sh

export PATH=/srv/git/tools:${PATH}

git-fetch-hourly /srv/git > /srv/git/git-fetch-hourly.log 2>&1
